# infra/supervisor/iievi.conf
; IIEVI process management — /etc/supervisor/conf.d/iievi.conf
; All processes run as the deploy user with Doppler-injected secrets.
; Logs go to /var/log/iievi/ (created by the hardening script) and are
; rotated daily by logrotate (see infra/logrotate/iievi).

[program:iievi-api]
command=doppler run -- /srv/iievi/apps/api/.venv/bin/gunicorn app.main:app
    --workers 4
    --worker-class uvicorn.workers.UvicornWorker
    --bind 127.0.0.1:8000
    --timeout 120
    --graceful-timeout 30
directory=/srv/iievi/apps/api
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/iievi/api.log
stderr_logfile=/var/log/iievi/api.error.log
environment=PATH="/srv/iievi/apps/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"

; Workers are split per queue group so each queue gets its spec'd concurrency:
; ai_conversations 8 | post_publishing 4 | creative_generation 3
; lead_outreach 4 | ad_management 2 | usage_tracking 4
; NOBODY consumes the dlq queue — it is drained by operations tooling only.
;
; Graceful shutdown: celery's default SIGTERM handling is a warm shutdown
; (finish the current task, then exit); stopwaitsecs=300 gives it up to five
; minutes before supervisor force-kills. This prevents the classic corruption:
; a post stuck at status=publishing because the worker died mid-publish.
; worker_max_tasks_per_child=1000 (set in celery config) recycles processes
; before slow memory leaks accumulate.

[program:iievi-worker-ai]
command=doppler run -- /srv/iievi/apps/api/.venv/bin/celery
    -A app.worker.celery_app:celery_app worker
    --loglevel=INFO
    --concurrency=8
    --hostname=ai@%%h
    -Q ai_conversations
directory=/srv/iievi/apps/api
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopwaitsecs=300
stdout_logfile=/var/log/iievi/worker-ai.log
stderr_logfile=/var/log/iievi/worker-ai.error.log
environment=PATH="/srv/iievi/apps/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[program:iievi-worker-publish]
command=doppler run -- /srv/iievi/apps/api/.venv/bin/celery
    -A app.worker.celery_app:celery_app worker
    --loglevel=INFO
    --concurrency=4
    --hostname=publish@%%h
    -Q post_publishing
directory=/srv/iievi/apps/api
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopwaitsecs=300
stdout_logfile=/var/log/iievi/worker-publish.log
stderr_logfile=/var/log/iievi/worker-publish.error.log
environment=PATH="/srv/iievi/apps/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[program:iievi-worker-creative]
command=doppler run -- /srv/iievi/apps/api/.venv/bin/celery
    -A app.worker.celery_app:celery_app worker
    --loglevel=INFO
    --concurrency=3
    --hostname=creative@%%h
    -Q creative_generation
directory=/srv/iievi/apps/api
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopwaitsecs=300
stdout_logfile=/var/log/iievi/worker-creative.log
stderr_logfile=/var/log/iievi/worker-creative.error.log
environment=PATH="/srv/iievi/apps/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[program:iievi-worker-ops]
command=doppler run -- /srv/iievi/apps/api/.venv/bin/celery
    -A app.worker.celery_app:celery_app worker
    --loglevel=INFO
    --concurrency=4
    --hostname=ops@%%h
    -Q lead_outreach,ad_management,usage_tracking
directory=/srv/iievi/apps/api
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopwaitsecs=300
stdout_logfile=/var/log/iievi/worker-ops.log
stderr_logfile=/var/log/iievi/worker-ops.error.log
environment=PATH="/srv/iievi/apps/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[program:iievi-beat]
command=doppler run -- /srv/iievi/apps/api/.venv/bin/celery
    -A app.worker.celery_app:celery_app beat
    --loglevel=INFO
directory=/srv/iievi/apps/api
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/iievi/beat.log
stderr_logfile=/var/log/iievi/beat.error.log
environment=PATH="/srv/iievi/apps/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[program:iievi-flower]
; Monitoring UI — bound to localhost only; access via SSH tunnel:
;   ssh -L 5555:127.0.0.1:5555 deploy@<vps>
command=doppler run -- /srv/iievi/apps/api/.venv/bin/celery
    -A app.worker.celery_app:celery_app flower
    --address=127.0.0.1
    --port=5555
directory=/srv/iievi/apps/api
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/iievi/flower.log
stderr_logfile=/var/log/iievi/flower.error.log
environment=PATH="/srv/iievi/apps/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[group:iievi]
programs=iievi-api,iievi-worker-ai,iievi-worker-publish,iievi-worker-creative,iievi-worker-ops,iievi-beat,iievi-flower
