# infra/supervisor/iievi.conf
; IIEVI process management — /etc/supervisor/conf.d/iievi.conf
; All processes run as the deploy user with Doppler-injected secrets.
; Logs go to /var/log/iievi/ (created by the hardening script) and are
; rotated daily by logrotate (see infra/logrotate/iievi).

[program:iievi-api]
command=doppler run -- /srv/iievi/apps/api/.venv/bin/gunicorn app.main:app
    --workers 4
    --worker-class uvicorn.workers.UvicornWorker
    --bind 127.0.0.1:8000
    --timeout 120
    --graceful-timeout 30
directory=/srv/iievi/apps/api
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/iievi/api.log
stderr_logfile=/var/log/iievi/api.error.log
environment=PATH="/srv/iievi/apps/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[program:iievi-worker]
command=doppler run -- /srv/iievi/apps/api/.venv/bin/celery
    -A app.worker.celery_app:celery_app worker
    --loglevel=INFO
    --concurrency=8
    -Q ai_conversations,post_publishing,creative_generation,lead_outreach,ad_management,usage_tracking
directory=/srv/iievi/apps/api
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopwaitsecs=60
stdout_logfile=/var/log/iievi/worker.log
stderr_logfile=/var/log/iievi/worker.error.log
environment=PATH="/srv/iievi/apps/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[program:iievi-beat]
command=doppler run -- /srv/iievi/apps/api/.venv/bin/celery
    -A app.worker.celery_app:celery_app beat
    --loglevel=INFO
directory=/srv/iievi/apps/api
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/iievi/beat.log
stderr_logfile=/var/log/iievi/beat.error.log
environment=PATH="/srv/iievi/apps/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[program:iievi-flower]
; Monitoring UI — bound to localhost only; access via SSH tunnel:
;   ssh -L 5555:127.0.0.1:5555 deploy@<vps>
command=doppler run -- /srv/iievi/apps/api/.venv/bin/celery
    -A app.worker.celery_app:celery_app flower
    --address=127.0.0.1
    --port=5555
directory=/srv/iievi/apps/api
user=deploy
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/iievi/flower.log
stderr_logfile=/var/log/iievi/flower.error.log
environment=PATH="/srv/iievi/apps/api/.venv/bin:/usr/local/bin:/usr/bin:/bin"

[group:iievi]
programs=iievi-api,iievi-worker,iievi-beat,iievi-flower
